# usershub/DEBIAN/postinst
#!/bin/sh
# postinst script for usershub
#
# see: dh_installdeb(1)

set -e
set -o xtrace

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

. /usr/share/debconf/confmodule

if [ -f /usr/share/debconf/confmodule ]; then
    db_version 2.0
    db_input high usershub/URL_APPLICATION || true
    db_input high usershub/PASS_METHOD || true
    db_input high usershub/FILL_MD5_PASS || true
    db_input low usershub/COOKIE_EXPIRATION || true
    db_input low usershub/ACTIVATE_API || true
    db_input low usershub/ACTIVATE_APP || true
    db_input high usershub/SECRET_KEY || true
    for p in GUN_NUM_WORKERS GUN_HOST GUN_PORT; do
        db_input high usershub/${p} || true
    done
    db_go
fi

for p in URL_APPLICATION PASS_METHOD FILL_MD5_PASS COOKIE_EXPIRATION ACTIVATE_API ACTIVATE_APP SECRET_KEY; do
    db_get usershub/${p}
    export ${p}="$RET"
done

# Generate config files
echo "Generate configuration" >&2
CONF_DIR=/etc/geonature
for f in "usershub.conf"; do
    envsubst <$CONF_DIR/$f.init >$CONF_DIR/$f
done

# Reload supervisor
if [ -f /etc/init.d/supervisor ]; then
    if which invoke-rc.d >/dev/null; then
        if invoke-rc.d supervisor status >/dev/null 2>&1; then
            invoke-rc.d supervisor reload 3>/dev/null || true
        fi
    else
        if /etc/init.d/supervisor status >/dev/null 2>&1; then
            /etc/init.d/supervisor reload 3>/dev/null || true
        fi
    fi
fi

case "$1" in
configure)

    echo "configure done"

    ;;

abort-upgrade | abort-remove | abort-deconfigure) ;;

*)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
