# deb_packages_files/usershub/{{cookiecutter}}usershub/etc/supervisor/conf.d/usershub-service.conf
[program:usershub2]
command = /usr/bin/usershub
autostart=true
autorestart=true
stdout_logfile = /var/log/usershub/errors_uhv2.log
redirect_stderr = true
