# deb_packages_files/taxhub/app/etc/supervisor/conf.d/taxhub-service.conf
[program:taxhub]
command = /usr/bin/taxhub
autostart=true
autorestart=true
stdout_logfile = /var/log/taxhub/errors_taxhub.log
redirect_stderr = true
